# files/logstash.target
[Unit]
Description=Logstash

[Install]
WantedBy=multi-user.target
